# cloudpan189/crontab
# 每小时执行一次
* */1 * * * /opt/cloudpan189-go/cloudpan189-go sign >> /var/log/cloudpan189_sign.log
